# jenkins.service
[Unit]
Description=Jenkins

[Service]
User=jenkins

WorkingDirectory=/home/jenkins
ExecStart=/usr/bin/java -jar /home/jenkins/jenkins.war

[Install]
WantedBy=multi-user.target
